# deploy/systemd/lexforge-api.service
[Unit]
Description=LexForge API (FastAPI / uvicorn)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=lexforge
Group=lexforge
WorkingDirectory=/opt/lexforge
EnvironmentFile=/opt/lexforge/.env
Environment=PYTHONPATH=/opt/lexforge
ExecStart=/opt/lexforge/.venv/bin/uvicorn apps.api.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
