# deploy/systemd/lexforge-web.service
[Unit]
Description=LexForge Web (Next.js)
After=network.target lexforge-api.service
Wants=lexforge-api.service

[Service]
Type=simple
User=lexforge
Group=lexforge
WorkingDirectory=/opt/lexforge/apps/web
EnvironmentFile=/opt/lexforge/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
